# src/.config/systemd/user/llama-swap.service
[Unit]
Description=local llama-swap for small models

[Service]
Type=simple
Restart=on-failure
RestartSec=5
StartLimitBurst=3
StartLimitInterval=30
OOMScoreAdj=900
ExecStart=/usr/bin/llama-swap -config %h/.config/llama-swap/config.yaml -listen 127.0.0.1:7748

[Install]
WantedBy=default.target
